# dist/lsws-prometheus-exporter.rc.gentoo.in
#!/sbin/runscript

depend()
{
	need net
	after sshd
}

start()
{
	ebegin "Starting LiteSpeed Web Server Prometheus Exporter"
	/usr/local/lsws-prometheus-exporter/lsws-prometheus-exporter %CERT_FILE% %KEY_FILE% &
}

